# myprojects/alertik_project/alertik_PHP/doc/mazanex.cron
* * * * * root /usr/bin/php5 /home/web/mazanex.com/private/leech.php >>/home/web/mazanex.com/logs/cron.txt ; /usr/bin/php5 /home/web/mazanex.com/private/leech-ime.php >>/home/web/mazanex.com/logs/cron-ime.txt ; /usr/bin/php5 /home/web/mazanex.com/private/calc.php >>/home/web/mazanex.com/logs/calc.txt ; /usr/bin/php5 /home/web/mazanex.com/private-arzlive/create-index.php >>/home/web/mazanex.com/logs/cron-index.txt  ; /usr/bin/php5 /home/web/mazanex.com/private/create-index.php >>/home/web/mazanex.com/logs/cron-index.txt ; /usr/bin/php5 /home/web/mazanex.com/private-arzlive/create-line-charts.php >>/home/web/mazanex.com/logs/cron-chart.txt ; /usr/bin/php5 /home/web/mazanex.com/private/create-line-charts.php >>/home/web/mazanex.com/logs/cron-chart.txt ; /usr/bin/php5 /home/web/mazanex.com/private/create-js-update.php >>/home/web/mazanex.com/logs/cron-js.txt ;
